# matrix/rootfs/etc/services.d/matrix/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Matrix
# Runs the Synapse daemon
# ==============================================================================
exec python3.6 -B -m synapse.app.homeserver -c /config/matrix.yaml
