# matrix/rootfs/etc/services.d/matrix/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Matrix
# Runs the Synapse daemon
# ==============================================================================
bashio::log.info "Starting Matrix Synapse..."
exec python3 -B -m synapse.app.homeserver -c /config/matrix.yaml
