# matrix/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Matrix
# Runs the NGINX daemon for Riot
# ==============================================================================
# Wait for Synapse to start
bashio::net.wait_for 8448

# Disable HA Authentication if front door is open
if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi

bashio::log.info "Starting NGINX..."
exec nginx
